# deployments/pawikan-sentinel.service
[Unit]
Description=Pawikan Sentinel - Sea Turtle Detection System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pawikan-sentinel
Environment="PATH=/home/pi/pawikan-sentinel/.venv/bin"
ExecStart=/home/pi/pawikan-sentinel/.venv/bin/uvicorn src.core.main:app --host 0.0.0.0 --port 8000 --no-reload
Restart=always
RestartSec=10
# Resource limits for Pi 4B
MemoryMax=1G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
